# supervisor/cron.conf
[program:cron]
command = cron -f -L 15
autorestart=true