# supervisor/user.conf
[program:user-startup]
command=sh /var/www/start.sh
directory=/var/www/
autorestart=false
startsecs=10