# supervisor/user.conf
[program:user-startup]
command=sh /var/www/start.sh
autorestart=false